# docker/services.d/h2o/run
#!/usr/bin/execlineb -P

/usr/local/bin/h2o -c /usr/local/etc/h2o.conf